# scripts/cluster-health/claude-cluster-health.service
[Unit]
Description=Nightly Claude homelab cluster-health digest -> Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Adjust REPO if your checkout lives elsewhere.
Environment=REPO=%h/projects/webgrip/homelab-cluster
ExecStart=%h/projects/webgrip/homelab-cluster/scripts/cluster-health/digest.sh
# Don't hammer if the LAN/cluster is unreachable; just log and move on.
TimeoutStartSec=900
