# scripts/cluster-health/claude-cluster-health.timer
[Unit]
Description=Run the Claude homelab cluster-health digest nightly

[Timer]
# 08:30 local time daily; Persistent catches up after the machine was off.
OnCalendar=*-*-* 08:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
